# systemd/gohome-once@bills.timer
[Unit]
Description=Timer for gohome bills

[Timer]
OnCalendar=00:02

[Install]
WantedBy=timers.target