# systemd/gohome-once@weather.timer
[Unit]
Description=Timer for gohome weather

[Timer]
OnCalendar=8:00

[Install]
WantedBy=timers.target